# systemd/macd-bot.timer
[Unit]
Description=MACD Bot daily trigger (Mon-Fri 3:55 PM ET)

[Timer]
OnCalendar=Mon-Fri *-*-* 15:55:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
