# systemd/macd-bot.service
[Unit]
Description=MACD Trading Bot
After=ibgateway.service
Requires=ibgateway.service

[Service]
Type=oneshot
WorkingDirectory=/home/jean/trading
EnvironmentFile=/home/jean/trading/.env
ExecStart=/home/jean/lean-env/bin/python /home/jean/trading/macd_bot.py
Restart=on-failure
RestartSec=60
TimeoutStartSec=180
